# puppet/modules/network_default/files/puppet-agent.service
[Unit]
Description=Puppet agent

[Service]
ExecStart=/opt/puppetlabs/bin/puppet agent

[Install]
WantedBy=multi-user.target
